# sv/pipewire-pulse/run
#!/bin/sh

sv start pipewire || exit 1
exec pipewire-pulse 2>&1
