# systemd/user/redshift.service.d/override.conf
[Install]
WantedBy=wm-session.target

[Unit]
After=wm-session.target

[Service]
Restart=always
RestartSec=2
