# sv/yambar-DP-1/run
#!/bin/sh

exec yambar --config="/home/gbrlsnchs/.config/yambar/ultrawide.yml" 2>&1
