# sv/yambar-primary/run
#!/bin/sh

exec yambar --config="$XDG_CONFIG_HOME/yambar/primary.yml" 2>&1
