# sv/rivertile/run
#!/bin/sh

exec rivertile -view-padding 6 -outer-padding 6 2>&1
