# sv/rivertile/run
#!/bin/sh

exec rivertile \
	-main-location right \
	-main-ratio 0.5 \
	-view-padding 1 \
	-outer-padding 16 \
	2>&1
