# sv/yambar-secondary/run
#!/bin/sh

exec yambar --config="$XDG_CONFIG_HOME/yambar/secondary.yml" 2>&1
