# systemd/user/wm-session.target
[Unit]
Description=WM session
Documentation=
BindsTo=graphical-session.target
Wants=graphical-session-pre.target
After=graphical-session-pre.target
