# services/xss-lock/run
#!/bin/sh

exec 2>&1
exec xss-lock -n /usr/libexec/xsecurelock/dimmer -l -- lockscreen
