# sv/wireplumber/run
#!/bin/sh

sv start "$XDG_DATA_HOME/service/pipewire" || exit 1
exec wireplumber 2>&1
