# sv/wireplumber/run
#!/bin/sh

sv start pipewire || exit 1
exec wireplumber 2>&1
