# sv/wob/run
#!/bin/sh

wobsock="$XDG_RUNTIME_DIR/wob.sock"

rm -f "$wobsock" && mkfifo "$wobsock"
exec tail -f "$wobsock" \
	| wob "#${FG_COLOR}ff" "#${BG_COLOR}ff" "#${FG_COLOR}ff" 2>&1
