# sv/yambar-HDMI-A-1/run
#!/bin/sh

exec yambar --config="$XDG_CONFIG_HOME/yambar/portable.yml" 2>&1
