# sv/yambar-HDMI-A-1/run
#!/bin/sh

exec yambar --config="/home/gbrlsnchs/.config/yambar/portable.yml" 2>&1
